# deploy/imprint-tunnel@.service
[Unit]
Description=Cloudflare Tunnel for Imprint
After=memo-clover.service

[Service]
Type=simple
User=%i
ExecStart=/usr/local/bin/cloudflared tunnel run my-tunnel
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
